# service/tor/run
#!/bin/sh
chown -Rv tord:tord /home/tord/

cat /etc/tor/torrc
cat /etc/os-release
echo 'nameserver 127.0.0.1' > /etc/resolv.conf

/usr/bin/tor -f /etc/tor/torrc
