# service/collectd/run
#!/bin/bash
set -euxo pipefail

exec collectd -f
